# cs261.conf
description "Store live trading data in SQLite"

start on runlevel [2345]
stop on runlevel [!2345]

# respawn if process quits unexpectedly
respawn

setuid tgdn
setgid www-data

script
/bin/bash <<EOT
	cd /home/tgdn/cs261
	python live_data_analysis.py
EOT
end script
